# systemd/user/teddy-http.service
[Unit]\nDescription=Teddy HTTP TTS endpoint (user)\nAfter=network.target\n\n[Service]\nType=simple\nWorkingDirectory=/home/gaf/teddy_bear_project_round2\nExecStart=/home/gaf/venv_teddy/bin/python3 /home/gaf/teddy_bear_project_round2/scripts/http_tts_server.py\nRestart=on-failure\nRestartSec=3\nEnvironment=PYTHONUNBUFFERED=1\n# Optional: override devices\n# Environment=TTS_USB_DEVICE=usbout\n# Environment=TTS_LOOPBACK_DEVICE=plughw:Loopback,0,0\n\n[Install]\nWantedBy=default.target\n